# profile.d/fastlane.sh
#!/usr/bin/env bash

# Weird fastlane issue
export PATH="$HOME/.fastlane/bin:$PATH"
