# profile.d/kubernetes.sh
#!/usr/bin/env bash

# This has saved me a lot of time...
alias k='kubectl'
